# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001     

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=32MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis-server]
command=redis-server
startretries = 10
autorestart = true

[program:httpd]
command=/usr/local/bin/toughrun httpd
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/toughwlan/httpd.log




